# service/deploy/systemd/grocerysearch-scraper.service
[Unit]
Description=Run GrocerySearch scraper once
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/grocerysearch/service
ExecStart=/opt/grocerysearch/service/deploy/run_job.sh scrape
